# svc/launchd/internal/testdata/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>{{.BundleID}}</string>
        <key>EnvironmentVariables</key>
        <dict>
            <key>GO_TEST_SERVER_ADDR</key>
            <string>{{.GoTestServerAddr}}</string>
        </dict>
        <key>ProgramArguments</key>
        <array>
            <!--Run-Remote-Tests -->
            <string>{{.GoTestBinary}}</string>
            <string>-test.count=1</string>
            <string>-test.run={{.GoTestName}}</string>
            <string>-test.timeout=30s</string>
            <string>-test.v=true</string>
        </array>
        <key>RunAtLoad</key>
        <true />
        <key>StandardErrorPath</key>
        <string>{{.StderrFile}}</string>
        <key>StandardOutPath</key>
        <string>{{.StdoutFile}}</string>
        <key>ExitTimeOut</key>
        <integer>60</integer>
        <!--Sockets-->
        <key>Sockets</key>
        <dict>
            <!--tcp-->
            <key>tcp</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.TCP}}</string>
                <key>SockType</key>
                <string>stream</string>
                <key>SockFamily</key>
                <string>IPv4</string>
            </dict>
            <!--udp-->
            <key>udp</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.UDP}}</string>
                <key>SockType</key>
                <string>dgram</string>
                <key>SockFamily</key>
                <string>IPv4</string>
            </dict>
            <!--tcp-multiple-->
            <key>tcp-multiple</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.TCPMultiple}}</string>
                <key>SockType</key>
                <string>stream</string>
            </dict>
            <!--udp-multiple-->
            <key>udp-multiple</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.UDPMultiple}}</string>
                <key>SockType</key>
                <string>dgram</string>
            </dict>
        </dict>
    </dict>
</plist>
